# edict.service
# ══════════════════════════════════════════════════════════════
# 三省六部 · systemd 服务
#
# 部署步骤:
#   sudo cp edict.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now edict
#
# 查看状态/日志:
#   sudo systemctl status edict
#   journalctl -u edict -f
# ══════════════════════════════════════════════════════════════
[Unit]
Description=三省六部 Multi-Agent Orchestration Platform
Wants=network-online.target redis.service postgresql.service
After=network-online.target redis.service postgresql.service

[Service]
Type=simple
User=edict
WorkingDirectory=/opt/edict
ExecStart=/opt/edict/edict.sh run
Restart=on-failure
RestartSec=5
KillMode=control-group

Environment=EDICT_DASHBOARD_HOST=0.0.0.0
Environment=EDICT_DASHBOARD_PORT=7891
Environment=EDICT_BACKEND_PORT=8000
Environment=PYTHONPATH=/opt/edict
Environment=PATH=/home/edict/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/edict/data /opt/edict/logs /opt/edict/.pids
PrivateTmp=true

[Install]
WantedBy=multi-user.target
